# pwa/deploy/tanksync.service
[Unit]
Description=TankSync PWA Server
Documentation=https://github.com/Techposts/LoRa-Water-Tank-Monitor
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=tanksync
Group=tanksync
WorkingDirectory=/home/tanksync/tanksync-pwa
ExecStart=/usr/bin/node --dns-result-order=ipv4first server/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/home/tanksync/tanksync-pwa/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/tanksync/tanksync-pwa/data
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tanksync

[Install]
WantedBy=multi-user.target
